# bookpoly-recorder.service
[Unit]
Description=BookPoly Market Data Recorder
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/bookpoly
Environment="PATH=/root/bookpoly/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/root/bookpoly/venv/bin/python /root/bookpoly/src/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Limites de recursos
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target

