# paper-trading.service
[Unit]
Description=BookPoly Paper Trading Bot
After=network.target
After=bookpoly-recorder.service
After=volatility-recorder.service
After=signals-recorder.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/bookpoly
Environment="PATH=/root/bookpoly/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/root/bookpoly/venv/bin/python -m bot.paper_trading --coins btc,eth,sol,xrp --verbose
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Limites de recursos
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target

