# scripts/bot15min.service
# Serviço 24/7 para o bot 15min Polymarket (BTC, ETH, SOL, XRP)
# Instalar: sudo cp scripts/bot15min.service /etc/systemd/system/
#          sudo systemctl daemon-reload
#          sudo systemctl enable bot15min
#          sudo systemctl start bot15min
#
# Comandos: status -> sudo systemctl status bot15min
#          logs   -> sudo journalctl -u bot15min -f
#          parar  -> sudo systemctl stop bot15min

[Unit]
Description=Bot 15min Polymarket (BTC, ETH, SOL, XRP)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/bookpoly
ExecStart=/root/bookpoly/.venv/bin/python scripts/bot_15min.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bot15min

# Encerramento gracioso (Ctrl+C cancela ordens abertas)
TimeoutStopSec=30
KillMode=mixed

[Install]
WantedBy=multi-user.target
