# scripts/bot4h.service
# Serviço 24/7 para o bot 4h Polymarket (BTC, ETH, SOL, XRP)
# Instalar: sudo cp scripts/bot4h.service /etc/systemd/system/
#          sudo systemctl daemon-reload
#          sudo systemctl enable bot4h
#          sudo systemctl start bot4h
#
# Comandos: status -> sudo systemctl status bot4h
#          logs   -> sudo journalctl -u bot4h -f
#          parar  -> sudo systemctl stop bot4h

[Unit]
Description=Bot 4h Polymarket (BTC, ETH, SOL, XRP)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/bookpoly
Environment=PYTHONUNBUFFERED=1
ExecStart=/root/bookpoly/venv/bin/python -u scripts/bot_4h.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bot4h

# Encerramento gracioso (Ctrl+C cancela ordens abertas)
TimeoutStopSec=30
KillMode=mixed

[Install]
WantedBy=multi-user.target

